# roles/provision/rexray/files/rexray.service.d/10-restart-policy.conf
[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=30s
